# files/etc/crontabs/root
59 1 * * 7 pihole updateGravity
@daily /usr/sbin/logrotate /etc/logrotate.conf
